# contrib/etc/systemd/system/tutuicmptunnel-server-reaper.service
[Unit]
Description=Reap tutuicmptunnel sessions

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'tuctl reaper'
